# skills/mlops/bert-text-classification/templates/supervisor.conf
; Supervisor config: keep BERT FastAPI microservice running
; Usage:
;   sudo cp bert-klassifisering.conf /etc/supervisor/conf.d/
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl status

[program:bert-klassifisering]
command=/path/to/venv/bin/python serve.py
directory=/path/to/model/dir
user=forge
autostart=true
autorestart=true
startretries=3
stopwaitsecs=10

environment=
    MODEL_DIR="./epost-model-final",
    HOST="127.0.0.1",
    PORT="8000",
    API_KEY=""

stdout_logfile=/home/forge/.supervisor/bert-klassifisering-out.log
stderr_logfile=/home/forge/.supervisor/bert-klassifisering-err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
